# mailnotify/crontab
*/20 * * * * root python /mailNotify.py
